# dot_config/systemd/user/restic-unlock.service
[Unit]
Description=Run restic unlock on repository during downtime

[Service]
Type=oneshot
ExecStart=/usr/bin/restic unlock
