# dot_config/systemd/user/restic-unlock.timer.tmpl
[Unit]
Description=Timer to run restic unlock regularly

[Timer]
OnCalendar=Mon..Fri 11,15:30:00
Persistent=false

[Install]
WantedBy=timers.target
