# dot_config/systemd/user/restic-backup.timer
[Unit]
Description=Timer for restic backup on home directory

[Timer]
OnCalendar=*:10:00
Persistent=true

[Install]
WantedBy=timers.target
